# generated/src/etc/inittab
fstab:123:boot:/sbin/fstab.lua
shell:123:respawn:/sbin/getty.lua
devtab:123:boot:/sbin/devtab.lua
